# runit/dictation/run
#!/bin/sh
exec chpst -u "$USER" "$HOME/.local/bin/dictation" serve
