# runit/lilium-palette/run
#!/bin/sh
exec chpst -u "$USER" "$HOME/.local/bin/lilium-palette" serve "$HOME/notes/notes.db"
