# openrc/lilium-palette
#!/sbin/openrc-run

name="lilium-palette"
description="Lilium Palette - Note link search"

command="$HOME/.local/bin/lilium-palette"
command_args="serve $HOME/notes/notes.db"
command_user="${USER}"
command_background=true
pidfile="/run/user/$(id -u)/lilium-palette.pid"

depend() {
    need localmount
    after bootmisc
}
